# docker/config/privoxy/supervisord-privoxy.ini
[program:privoxy]
command=privoxy-start.sh
numprocs=1
autostart=true
autorestart=true
