# docker/openvpn/supervisord-openvpn.ini
[program:openvpn]
command=openvpn-start
#stderr_logfile=/var/log/supervisor/openvpn-stderr.log
#stdout_logfile=/var/log/supervisor/openvpn-stdout.log
numprocs=1
autostart=true
autorestart=true
