# docker/privoxy/supervisord-privoxy.ini
[program:privoxy]
command=privoxy-start.sh
#stderr_logfile=/var/log/supervisor/privoxy-stderr.log
#stdout_logfile=/var/log/supervisor/privoxy-stdout.log
numprocs=1
autostart=true
autorestart=true
