# docker/config/openvpn/supervisord-openvpn.ini
[program:openvpn]
command=openvpn-start.sh
numprocs=1
autostart=true
autorestart=true
